# infra/microserver/mihomo-container-proxy.socket
[Unit]
Description=Mihomo container proxy socket (172.18.0.1:17890 -> 127.0.0.1:7890)
Requires=mihomo.service
After=mihomo.service

[Socket]
ListenStream=172.18.0.1:17890
NoDelay=true

[Install]
WantedBy=sockets.target
